# mayan/files/supervisor.conf
{% from "mayan/map.jinja" import server with context %}
{% from "mayan/map.jinja" import api with context %}

{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}

[program:mayan-web]
command={{ server.dir.base }}/bin/gunicorn_start
directory = {{ server.dir.base }}/site
user = root
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile=/var/log/mayan-edms/mayan_app.log 
stderr_logfile=/var/log/mayan-edms/mayan_error.log

#[program:mayan-worker]
#command = {{ server.dir.base }}/bin/python {{ server.dir.base }}/bin/mayan-edms.py celery --settings=mayan.settings.production worker -Ofair -l ERROR
#directory = {{ server.dir.base }}/site
#user = www-data
#stdout_logfile = /var/log/mayan-edms/worker-stdout.log
#stderr_logfile = /var/log/mayan-edms/worker-stderr.log
#autostart = true
#autorestart = true
#startsecs = 10
#stopwaitsecs = 10
#killasgroup = true
#priority = 998
#
#[program:mayan-beat]
#command = {{ server.dir.base }}/bin/python {{ server.dir.base }}/bin/mayan-edms.py celery --settings=mayan.settings.production beat -l ERROR
#directory = {{ server.dir.base }}/site
#user = www-data
#numprocs = 1
#stdout_logfile = /var/log/mayan-edms/beat-stdout.log
#stderr_logfile = /var/log/mayan-edms/beat-stderr.log
#autostart = true
#autorestart = true
#startsecs = 10
#stopwaitsecs = 1
#killasgroup = true
#priority = 998

{%- if service.name == "api" and service.type == 'mayan' %}

[program:mayan-api]
command=python /srv/mayan/site/manage.py run_api
autostart = true
autorestart=true
user=mayan
{%- endif %}