# mayan/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}

[program:mayan_server]

command=/srv/mayan/bin/gunicorn_start

stdout_logfile=/srv/mayan/logs/mayan_app.log 
stderr_logfile=/srv/mayan/logs/mayan_error.log

user=mayan
autostart=true
autorestart=true

{#
[program:pyro_api]
command=python /srv/mayan/site/manage.py run_api
autostart = true
autorestart=true
user=mayan
#}